# server-next/deploy/all-in-one/supervisord.conf
; Takumi server-next — single container (supervisord)
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:postgres]
command=/usr/local/bin/postgres -D %(ENV_PGDATA)s
user=postgres
autostart=true
autorestart=true
startretries=20
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:takumi-legacy-login]
command=/usr/local/bin/takumi-start-legacy.sh
autostart=true
autorestart=unexpected
startretries=5
startsecs=5
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:takumi-game-host]
command=/opt/takumi/game-host/Takumi.Server.GameHost
directory=/opt/takumi/game-host
autostart=true
autorestart=unexpected
startretries=5
startsecs=3
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
